# packages/project-template/base/infra/supervisor.example.conf
[program:harmeese-website]
command=/bin/bash -lc 'cd %(ENV_PROJECT_DIR)s && WEBSITE_PORT=%(ENV_WEBSITE_PORT)s pnpm start'
autostart=true
autorestart=true
stdout_logfile=/tmp/harmeese-website.log
stderr_logfile=/tmp/harmeese-website.err

[program:harmeese-agent]
command=/bin/bash -lc 'cd %(ENV_PROJECT_DIR)s && ./infra/start-agent.sh'
autostart=true
autorestart=true
stdout_logfile=/tmp/harmeese-agent.log
stderr_logfile=/tmp/harmeese-agent.err
